# syncthing/rootfs/etc/services.d/synct/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Syncthing
# ==============================================================================
bashio::log.info "Starting Syncthing..."

exec syncthing --config=/config/syncthing/ --data=/config/syncthing/