# jenkins/rootfs/etc/s6-overlay/s6-rc.d/jenkins/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: Jenkins
# Runs Jenkins
# ==============================================================================
declare -a options
bashio::log.info 'Init Jenkins...'
if [ ! -d /config/jenkins  ]; then
	mkdir -p /config/jenkins 
fi
export JENKINS_HOME="/config/jenkins"
#export JENKINS_PREFIX="/hassio/ingress/973b0e7c-zhomejenkins"
#export JENKINS_PREFIX="/jenkins"
export HOME="/root"
#export JENKINS_PREFIX="/hassio/ingress/local_zhomejenkins"
#export PATH=/root/.arduino15/packages/esp8266/tools:${PATH}
chmod 755 -R /root/.arduino15/packages/esp8266/tools
readonly arduino_cli_version=`arduino-cli version`

bashio::log.info "Version Arduino : ${arduino_cli_version}"
bashio::log.info "Version JVM : " 
readonly jvm_version=`java -version`
#arduino-cli core update-index 
#arduino-cli core install esp8266:esp8266 

bashio::log.info 'Starting Jenkins...'

# Run Jenkins
exec jenkins
#/usr/local/bin/jenkins.sh