# jellyfin/rootfs/etc/s6-overlay/s6-rc.d/jellyfin-init/run
#!/usr/bin/with-contenv bashio

if [ ! -d /config/jellyfin/logs/nginx ]; then
	mkdir -p /config/jellyfin/logs/nginx
	touch /config/jellyfin/logs/nginx/access.log
	touch /config/jellyfin/logs/nginx/error.log 
fi 
if [ ! -d /config/jellyfin ]; then
	mkdir -p /config/jellyfin
fi
if [ ! -d /config/jellyfin/config ]; then
	mkdir -p /config/jellyfin/config
fi
if [ ! -d /config/jellyfin/cache ]; then
	mkdir -p /config/jellyfin/cache
fi
lshw -C display
ls -l /dev/dri
addgroup jellyfin video

#ffmpeg -v verbose -init_hw_device vaapi=va:/dev/dri/renderD128 -init_hw_device opencl@va
