# jellyfin/rootfs/etc/s6-overlay/s6-rc.d/jellyfin/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Jellyfin
# Starts Jellyfin
# ==============================================================================

bashio::log.info "Starting Jellyfin..."

exec jellyfin 
