# jenkins/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: jenkins
# Configures NGINX
# ==============================================================================
declare -a options
bashio::log.info 'Init Nginx...'
if [ ! -d /config/jenkins/logs/nginx ]; then
	mkdir -p /config/jenkins/logs/nginx
	touch /config/jenkins/logs/nginx/access.log
	touch /config/jenkins/logs/nginx/error.log 
fi 
whoami
tmp=`whoami`
cat /etc/passwd 
id "${tmp}"
chmod -R 755 /config/jenkins 
usermod -aG jenkins www-data
usermod -aG jenkins nginx